# deploy/supervisor/mail-scheduler-worker.conf
[program:mail-scheduler-worker]
command=/usr/bin/python3 -m src.services.worker
directory=/opt/mail-analyzer
autostart=true
autorestart=true
startsecs=5
startretries=10
stopsignal=INT
stopasgroup=true
killasgroup=true
environment=PYTHONPATH="/opt/mail-analyzer"
stdout_logfile=/opt/mail-analyzer/logs/worker.out.log
stderr_logfile=/opt/mail-analyzer/logs/worker.err.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5
